# roles/dyndns/templates/dyndns.timer.j2
[Unit]
Description=update dynamic DNS once an hour once online
Wants=network-online.target
After=network-online.target

[Timer]
OnBootSec=5
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target