# roles/dyndns/templates/dyndns.service.j2
[Unit]
Description=dynamic DNS

[Service]
Type=oneshot
ExecStart=/home/dyndns/dyndns
User=dyndns
Group=systemd-journal